# debian/nginx.postrm.in
#!/bin/sh -e
D=/dev/null

d_s_h() {
  if which deb-systemd-helper >$D; then
    deb-systemd-helper "$@"
  fi
}

case "$1" in
purge)
  rm -rf '/%path_data%' '/%path_log%' '/%root_path_data%' '/%root_path_log%'

  for i in nginx root-nginx; do s="$i.service"
    update-rc.d "$i" remove >$D
    d_s_h purge "$s" >$D
    d_s_h unmask "$s" >$D
  done

  if [ -d /run/systemd/system ] ; then
    systemctl --system daemon-reload >$D || true
  fi
;;
remove)
  for i in nginx root-nginx; do s="$i.service"
    d_s_h mask "$s" >$D
  done
;;
upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)
;;
*)
  echo "postrm called with unknown argument \`$1'" >&2
  exit 1
;;
esac

exit 0
