# debian/nginx.root-nginx.init.in
#!/bin/sh

### BEGIN INIT INFO
# Provides:          root-nginx
# Required-Start:    $local_fs $remote_fs $network $syslog $named
# Required-Stop:     $local_fs $remote_fs $network $syslog $named
# Default-Start:     
# Default-Stop:      0 1 6
# Short-Description: starts the nginx web server with root privileges
# Description:       starts nginx using start-stop-daemon
### END INIT INFO

NAME='root-nginx'
DESC='nginx - high performance web/proxy server (ROOT INSTANCE)'

. '/%path_prefix%/init-functions'
