# debian/nginx.postrm.in
#!/bin/sh
set -e
D=/dev/null

d_s_h() {
  if which deb-systemd-helper >$D; then
    deb-systemd-helper "$@"
  fi
}

case "$1" in
purge)
  rm -rf '/%path_data%' '/%path_log%'

  s="nginx.service"
  update-rc.d "$i" remove >$D
  d_s_h purge "$s" >$D
  d_s_h unmask "$s" >$D

  if [ -d /run/systemd/system ] ; then
    systemctl --system daemon-reload >$D || true
  fi
;;
remove)
  s="nginx.service"
  d_s_h mask "$s" >$D
;;
upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)
;;
*)
  echo "postrm called with unknown argument \`$1'" >&2
  exit 1
;;
esac

exit 0
