# debian/nginx-bin.postinst.in
#!/bin/sh
set -e
D=/dev/null

case "$1" in
triggered)
  pidof '%file_binary%' >$D || exit 0

  pid_normal='-1'
  if [ -s '%file_pid%' ]; then
    echo "Triggering reload: nginx ..."
    if [ -x "/etc/init.d/nginx" ]; then
      invoke-rc.d nginx reload || true
    else
      pid_normal=$(cat '%file_pid%' 2>$D || true)
      kill -HUP ${pid_normal} || true
    fi
  fi

  pid_others=$(pidof '%file_binary%' | xargs -r ps -o sid= | \
    sort -u | grep -Ev "^(|${pid_normal})\$" || true)

  [ -n "${pid_others}" ] || exit 0

  echo "Triggering reload: other instances (${pid_others}) ..."
  kill -HUP ${pid_others} || true

  exit 0
;;
abort-upgrade|abort-remove|abort-deconfigure|configure)
;;
*)
  echo "postinst called with unknown argument \`$1'" >&2
  exit 1
;;
esac

pidof '%file_binary%' >$D || exit 0

if [ -s '%file_pid%' -a -x "/etc/init.d/nginx" ]; then
  if ! invoke-rc.d nginx upgrade; then
    invoke-rc.d nginx restart || true
  fi
fi

exit 0
