# debian/nginx.postinst.in
#!/bin/sh -e
D=/dev/null

d_s_h() {
  if which deb-systemd-helper >$D; then
    deb-systemd-helper "$@"
  fi
}

setup_access_int() {
  chown $1 "$4"
  chgrp $2 "$4"
  chmod $3 "$4"
}

setup_access() {
  ## allow local admin to override
  stat=$(dpkg-statoverride --list "$4" || true)
  if [ -n "${stat}" ]; then
    setup_access_int ${stat}
  else
    setup_access_int "$@"
  fi
}

setup_dir() {
  if [ ! -e "$4" ]; then
    mkdir -p "$4"
  fi
  setup_access "$@"
}

setup_file() {
  if [ ! -e "$4" ]; then
    touch "$4"
  fi
  setup_access "$@"
}

case "$1" in
configure)
  for i in '/%path_log%' '/%root_path_log%'; do
    setup_dir root '%nginx_log_group%' 0755 "$i"
  done

  setup_dir '%nginx_user%' root 0750 '/%path_data%'
  setup_dir root root 0750 '/%root_path_data%'

  ## secure default logfiles on fresh installations
  if [ -z "$2" ]; then
    for i in '%file_log_access%' '%file_log_error%'; do
      setup_file '%nginx_user%' '%nginx_log_group%' 0640 "$i"
    done

    for i in '%root_file_log_access%' '%root_file_log_error%'; do
      setup_file root '%nginx_log_group%' 0640 "$i"
    done
  fi
;;
abort-upgrade|abort-remove|abort-deconfigure)
;;
*)
  echo "postinst called with unknown argument \`$1'" >&2
  exit 1
;;
esac

for i in nginx root-nginx; do s="$i.service"
  d_s_h unmask "$s" >$D || true
  if d_s_h --quiet was-enabled "$s"; then
    d_s_h enable "$s" >$D || true
  else
    d_s_h update-state "$s" >$D || true
  fi
done

case "$1" in
configure|abort-upgrade)
  for i in nginx root-nginx; do
    [ -x "/etc/init.d/$i" ] || continue
    update-rc.d "$i" defaults || exit $?
    invoke-rc.d "$i" try-restart || exit $?
  done
;;
esac

exit 0
